# etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
user=root
autorestart=true
priority=400