# etc/supervisor/conf.d/xrdp.conf
[program:xrdp]
command=//usr/local/sbin/xrdp --nodaemon
user=root
autorestart=true
priority=400