# etc/supervisor/conf.d/xrdp-sesman.conf
[program:xrdp-sesman]
command=/usr/local/sbin/xrdp-sesman --nodaemon
user=root
autorestart=true
priority=400